# 80-restart-services-on-wlan0-conn.rules
ATTRS{idVendor}=="7392", ATTRS{idProduct}=="7811", RUN+="/etc/restart-services.sh"
