# 85-restart-services-on-wlan1-conn.rules
ATTRS{idVendor}=="148f", ATTRS{idProduct}=="7601", RUN+="/etc/restart-services.sh"
