# flask_app.conf
[program:flask_app]
command = flask run -h 10.5.5.1
directory = /var/www/wifiextender
autostart = true
autorestart = true
